# conf/devlog-tasks.target
[Unit]
Description=Devlog RQ worker pool

[Install]
WantedBy=devlog.target
