# conf/devlog-task@.service
[Unit]
Description=Devlog task worker %I
After=network.target

[Service]
User=username
Group=www-data
WorkingDirectory=/home/username/devlog
Environment="FLASK_ENV=production"
Environment="SENTRY_PUBKEY=your-sentry-pub-key"
Environment="SENTRY_PROJECT=your-sentry-project-id"
Environment="SQLALCHEMY_DATABASE_URI=postgresql://username:password@localhost:5432/dbname"
Environment="SECRET_KEY=your-application-secret-key"
Environment="REDIS_URL=redis://"
ExecStart=/home/username/devlog/venv/bin/rq worker devlog-tasks
Restart=always

[Install]
WantedBy=devlog-tasks.target
