# conf/devlog.target
[Unit]
Description=Devlog things

[Install]
WantedBy=multi-user.target
